# conf/ContainerCleaner.cron
#Enable ContainerCleaner
00 02 * * * root [[ `hostname` == YOUR-HOSTNAME ]]; [[ $? == 0 ]] && python3 /usr/share/ContainerCleaner/src/ContainerCleaner.py > /dev/null 2>&1

