# conf/ContainerCleaner.cron
#Enable ContainerCleaner
00 02 * * * /usr/bin/python3 /usr/share/ContainerCleaner/src/ContainerCleaner.py

